# .config/systemd/user/borg-backup.timer
[Unit]
Description=Runs Borg Backup Service

[Timer]
OnCalendar=*:0/10
OnBootSec=5min
Unit=borg-backup.service

[Install]
WantedBy=timers.target
