# .config/systemd/user/ipo.service
[Unit]
Description=Automatic IPO

[Service]
Type=simple
ExecStart=/home/kayc/Code/Python/MeroShare-IPO/ipo.py --default True


