# .config/systemd/user/safeeyes.timer
[Unit]
Description=Safeeyes Timer

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target
