# .config/systemd/user/wallpaper.service
[Unit]
Description=Automatic Wallpaper

[Service]
Type=simple
ExecStart=/home/kayc/.dotfiles/.local/share/bin/autowallpaper.sh
