# .config/systemd/user/wallpaper.timer
[Unit]
Description=Automatic Wallpaper Change Timer

[Timer]
OnBootSec=10

[Install]
WantedBy=timers.target
