# .config/systemd/user/wallpaper.timer
[Unit]
Description=Automatic Wallpaper Change Timer

[Timer]
OnBootSec=100
OnCalendar=07:00
Persistent=true

[Install]
WantedBy=timers.target
